# crontab
* * * * * cd /app && vendor/bin/crunz schedule:run
